# ai.openclaw.console.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.console</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/opt/node/bin/node</string>
      <string>/Users/huangzhenfeng/WorkSpace/openclaw-watchdog/console/server.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/huangzhenfeng/WorkSpace/openclaw-watchdog</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/opt/node/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>/Users/huangzhenfeng</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/huangzhenfeng/.openclaw/logs/console.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/huangzhenfeng/.openclaw/logs/console.err.log</string>
  </dict>
</plist>
